# WN-AX2033GR-2.01.000/package/network/services/pptpd/files/etc/hotplug.d/iface/70-change-lanip
#!/bin/sh

if [ "$INTERFACE" = lan ] ; then
	if [ "$ACTION" = ifup ]; then
		enabled=$(uci get pptpd.pptpd.enabled)
		[ -z "$enabled" -o "$enabled" = "0" ] && return 0
		/etc/init.d/pptpd restart
	fi
fi
