# WN-AX2033GR-2.01.000/package/ramips/ui/luci-mtk/src/themes/iodata/root/etc/uci-defaults/luci-theme-iodata
#!/bin/sh
uci batch <<-EOF
	set luci.themes.OpenWrt=/luci-static/iodata
        commit luci
EOF
