# WN-AX2033GR-2.01.000/package/base-files/files/etc/hotplug.d/iface/10-mstc_ethernet_lamponoff
[ ifup = "$ACTION" ] && {
	[ -f /etc/leds.d/lamp_on_off.sh ] && {
		/etc/leds.d/lamp_on_off.sh &
	}	
}


