# WN-AX2033GR-2.01.000/package/ramips/applications/mstc_eee/files/mstc_eee.init
#!/bin/sh /etc/rc.common

load_eee() {
	local enable
  config_get enable $1 enable

  [ -n "$enable" -o "$enable" = "1" -o "$enable" = "0" ] && {
    echo "Start switch EEE (smart_speed) "
    local command="/usr/sbin/switch-eee.sh $enable"
    eval $command
  }
}

setup_eee() {
	config_load system
	config_foreach load_eee eee
}

