# WN-AX2033GR-2.01.000/package/kernel/mtd-locked/files/mstc_mtd_locked
#!/bin/sh /etc/rc.common

START=19
DEBUGFLAG=`sys debugflag | sed 's/DebugFlag: //g'`
start() {
	if [ $DEBUGFLAG != 1 ]; then
		insmod mtd_locked
	fi
}
